# wa-approve/wa-approve.service
[Unit]
Description=WA approve/deselect surface for follow-up WABA sends (Flask via gunicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/wa
# Port 8101 matches WA_APPROVE_PORT in /root/wa/wa_approve.env (verified 11 Jul 2026).
# The app self-loads wa_approve.env from its own folder, so no EnvironmentFile= needed.
# --timeout 300: one LIVE batch can fire up to WA_DAILY_CAP sends in a single POST;
# a short worker timeout would kill it mid-batch.
ExecStart=/root/wa/venv/bin/gunicorn -w 1 -b 127.0.0.1:8101 --timeout 300 wa_approve:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
