# call-hook/call-hook.service
[Unit]
Description=MyOperator call-webhook capture (Flask via gunicorn) - Phase A, passive
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/wa/call-hook
ExecStart=/root/wa/venv/bin/gunicorn -w 1 -b 127.0.0.1:8098 call_hook_capture:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
